# tank-highscore-bot.service
[Unit]
Description=Discord Tank Highscore Bot
After=network.target

[Service]
Type=simple
User=botuser
WorkingDirectory=/opt/tank_highscore_bot
ExecStart=/opt/tank_highscore_bot/venv/bin/python /opt/tank_highscore_bot/bot.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

# Optional: point to env file if you prefer systemd-managed env
# EnvironmentFile=/opt/tank_highscore_bot/.env

[Install]
WantedBy=multi-user.target